# library.timer
[Unit]
Description=Timer for library check

[Timer]
OnCalendar=*-*-* 2:00:00
Persistent=true

[Install]
WantedBy=timers.target